# systemd-demo/debian/systemd-demo.service
[Unit]
Description=Simple demonstration service

[Service]
ExecStart=/bin/ping google.com

[Install]
WantedBy=multi-user.target
